# claude-code-webui/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[include]
files = /etc/supervisor/conf.d/*.conf

[program:nginx]
command=nginx -c /opt/nginx/nginx.conf -g "daemon off;"
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:claude-code-webui]
command=/usr/local/bin/claude-code-webui --host 127.0.0.1 --port 3000
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[eventlistener:stop_on_any_exit]
command=bash -c "printf 'READY\n'; while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
