# ntopng/supervisord.conf
[supervisord]
nodaemon=true

#[program:redis]
#command=/usr/bin/redis-server

[program:ntopng]
command=/usr/bin/ntopng -c /etc/ntopng/ntopng.conf
